# .config/systemd/user/update-checker.timer
[Unit]
Description=Timer for update checker service

[Timer]
OnCalendar=hourly
Persistent=True

[Install]
WantedBy=timers.target
